# cron/autowatch.cron
*/2 * * * * /bin/bash /opt/autowatch/bin/monitor.sh
